# xserver/init/startup/system-controller/raft-startup.service
[Unit]
Description=raft-startup

[Service]
ExecStart=/usr/bin/raft-startup
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
